# docker/cron/crontab
2 * * * * root /usr/local/bin/php /usr/local/apache2/htdocs/admin/cli/cron.php >/dev/null
# An empty line is required at the end of this file for a valid cron file.
